# setup/pizza.service
[Unit]
Description=PizzaPi Startup Service

[Service]
ExecStart=/bin/bash /root/PizzaPi/setup/boot.sh

[Install]
WantedBy=default.target
