# setup/pizza.service
[Unit]
Description=PizzaPi Startup Service

[Service]
ExecStart=~/PizzaPi/setup/boot.sh

[Install]
WantedBy=default.target
